# systemd/check-megaraid-monthly.service
[Unit]
Description=MegaRAID monthly health summary email
Documentation=file:///usr/local/share/doc/check-megaraid/README.md

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/check-megaraid --monthly
Nice=10
